# packaging/js/scripts/control/deb/postrm
#!/bin/sh

update-rc.d -f ${pkg.name} remove
