# msa/js-executor/src/main/scripts/control/deb/postinst
#!/bin/sh

chown -R ${pkg.user}: ${pkg.logFolder}
chown -R ${pkg.user}: ${pkg.installFolder}
update-rc.d ${pkg.name} defaults

