# transport/coap/src/main/scripts/control/deb/postinst
#!/bin/sh

chown -R ${pkg.name}: ${pkg.logFolder}
chown -R ${pkg.name}: ${pkg.installFolder}
update-rc.d ${pkg.name} defaults

