# transport/coap/src/main/scripts/control/rpm/prerm
#!/bin/sh

if [ $1 -eq 0 ] ; then
        # Package removal, not upgrade
        systemctl --no-reload disable --now ${pkg.name}.service > /dev/null 2>&1 || :
fi
